# scripts/systemd/pi-stack-deploy.service
# Runs scripts/deploy.sh once. Triggered by pi-stack-deploy.timer.
# Adjust User and WorkingDirectory to where you cloned the repo on the Pi.
[Unit]
Description=Deploy pi-stack from git
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/pi-stack
ExecStart=/home/pi/pi-stack/scripts/deploy.sh
# Don't hammer the Docker socket / registry if a run overruns the interval.
TimeoutStartSec=600
